# roles/netcat/templates/nc-msg.service.j2
[Unit]
Description=nc-msg
After=Network.target

[Service]
Type=simple
ExecStart={{ netcat_file_path }}/nc-msg.sh
User={{ netcat_user }}
Group={{ netcat_group }}

[Install]
WantedBy=default.target